# supervisor/aria2chrome.ini
[program:aria2chrome]
command=bash -c ". %(ENV_HOME)s/.bashrc; aria2c --conf-path=/Users/ben/.aria2/aria2-chrome.conf"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=%(ENV_HOME)s/.aria2/aria2-chrome-log.txt
stderr_logfile=%(ENV_HOME)s/.aria2/aria2-chrome-error.txt
user=ben
